# files/standard/pi/systemd/system/attach-pijuice-ups.service
[Unit]
Description=attach-pijuice UPS
Wants=untrustedhost-upsd-config.service
Before=untrustedhost-upsd-config.service

[Service]
Type=oneshot
ExecStart=/usr/lib/untrustedhost/scripts/attach-pijuice.sh
