# files/standard/systemd/system/incron.service.d/10-fitstat.conf
[Service]
ExecStartPre=-/usr/bin/mkdir -p /dev/fitusb
ExecStartPre=-/usr/bin/mkdir -p /run/untrustedhost/
ExecStartPre=-/usr/bin/touch /run/untrustedhost/sys_status
