# files/standard/systemd/system/attach-hidups@.service
[Unit]
Description=attach-hidups for %I

[Service]
Type=oneshot
ExecStart=/usr/bin/echo %i
