# files/standard/systemd/system/tmux-console.service
[Unit]
Description=tmux console
[Service]
Type=oneshot
ExecStart=/usr/bin/tmux -L console new-session -d -s console bash
RemainAfterExit=yes
