# files/standard/pi/systemd/system/pijuice.service.d/10-nut-driver.conf
[Unit]
Wants=attach-pijuice-ups.service
