# files/lite/systemd/system/ssh.service.d/ssh-keygen.conf
[Unit]
After=sshd-keygen.target
Wants=sshd-keygen.target
