# files/standard/systemd/system/nut-driver.service.d/10-sleep.conf
[Service]
# give the ups drivers time to settle
ExecStartPort=/bin/sleep 1
