# files/standard/systemd/system/networking.service.d/10-create-interface-cfg.conf
[Unit]
After=imd.service systemd-networkd.service
Requires=imd.service
