# files/standard/systemd/system/nut-server.service.d/10-start-aliases.conf
[Service]
ExecStartPost=/usr/lib/untrustedhost/scripts/start-upsd-aliases.sh
